# anytls-panel.service
[Unit]
Description=AnyTLS Panel - 节点统一管理面板
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/anytls-panel
ExecStart=/opt/anytls-panel/venv/bin/gunicorn -w 2 -b 0.0.0.0:8866 --timeout 60 app:app
Restart=always
RestartSec=5
Environment=PYTHONUNBUFFERED=1
Environment=ANYTLS_TRAFFIC_API_TOKEN_FILE=/opt/anytls-panel/.traffic_api_token

[Install]
WantedBy=multi-user.target
